# etc/crontab/backend
# Edit this file to introduce tasks to be run by cron.
# 
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
# 
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').# 
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
# 
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
# 
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command

# STAGE ONE
# BACKPAN & CPAN updates
46 * * * * bash /media/web/pause/pause-run.sh        >/media/web/pause/pause-run.out     2>&1
50 * * * * bash /media/web/pause/backpan.sh         >>/media/web/pause/logs/backpan.out  2>&1

# STAGE TWO
# Uploads tasks
*/5         * * * * bash /media/backend/cpantesters/autorun-data1.sh >>/media/backend/cpantesters/cron/autorun-data1.out 2>&1

# STAGE THREE
# Generate tasks
0,15,30,45  0-23 * * * bash /media/backend/cpantesters/autorun-data3.sh >>/media/backend/cpantesters/cron/autorun-data3.out 2>&1
45    0-23 * * * bash /media/backend/cpantesters/autorun-data4.sh >>/media/backend/cpantesters/cron/autorun-data4.out 2>&1
0     8    * * 6 bash /media/backend/cpantesters/autorun-data2.sh >>/media/backend/cpantesters/cron/autorun-data2.out 2>&1

# STAGE FOUR
# Website Monitors
0,15,30,45  * * * * bash /var/www/reports/toolkit/reports-monitor.sh >>/var/www/reports/toolkit/logs/monitor-cron.log  2>&1
0           * * * * perl /var/www/reports/toolkit/reports-recent.sh  >>/var/www/reports/toolkit/logs/recent-cron.log   2>&1

# Report Site Builders
#x#1,16,31,46  0-23 * * * bash /var/www/reports/toolkit/reports-builder-distro.sh --start >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1
#x#6,21,36,51 12-23 * * * bash /var/www/reports/toolkit/reports-builder-author.sh --start >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1
#x#1 0              * * * bash /var/www/reports/toolkit/reports-builder-author.sh --stop  >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1
# 2015-10-07 - disabled the author stop request, as the new server looks to be able to keep up. Barbie.

#-- don't run unless needed. distros build all the time, authors build in the morning only (helps to manage resources)
#0 0,12,21       * * * bash /var/www/reports/toolkit/reports-builder-distro.sh --stop  >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1

# STAGE FIVE
# Website updates
40  4 * * * bash /media/backend/cpantesters/autorun-web1.sh   >>/media/backend/cpantesters/cron/autorun-web1.out 2>&1
59  0 * * * bash /media/backend/cpantesters/autorun-web2.sh   >>/media/backend/cpantesters/cron/autorun-web2.out 2>&1
*/5 * * * * bash /media/backend/cpantesters/autorun-web3.sh   >>/media/backend/cpantesters/cron/autorun-web3.out 2>&1
0   * * * * bash /var/www/reports/toolkit/sponsors.sh         >>/var/www/reports/toolkit/sponsors.out 2>&1

# STAGE SIX
# Mailers
15 2          * * *  bash /media/backend/cpantesters/autorun-mail1.sh >>/media/backend/cpantesters/cron/autorun-mail1.out 2>&1
# -- note uploads mailer is in autorun-back2.sh

# STAGE SEVEN
# Backups
32 0          * * 0  bash /media/backend/cpantesters/autorun-back1.sh >>/media/backend/cpantesters/cron/autorun-back1.out 2>&1
31 2          * * *  bash /media/backend/cpantesters/autorun-back2.sh >>/media/backend/cpantesters/cron/autorun-back2.out 2>&1
35 0          * * *  perl /var/www/dbdump.pl >>/var/www/backups/dbdump.cron 2>&1


# ---- old entries, here for reference only

# Release Updates
#-- these have been disabled as the release is no longer maintained
#disabled#10 4,10,16,22 * * * bash /var/www/reports/toolkit/reports-release.sh >>/var/www/reports/toolkit/logs/release-cron.log 2>&1
#disabled#10 0          * * * bash /media/backend/cpantesters/autorun-back3.sh >>/media/backend/cpantesters/cron/autorun-back3.out 2>&1

# Old Uploads entries
#10,40  0-23 * * * bash /media/backend/cpantesters/autorun-data1.sh >>/media/backend/cpantesters/cron/autorun-data1.out 2>&1
#20,30,50  6 * * * bash /media/backend/cpantesters/autorun-data1.sh >>/media/backend/cpantesters/cron/autorun-data1.out 2>&1

# Old Generate tasks
#30   0-23 * * * bash /media/backend/cpantesters/autorun-data3.sh >>/media/backend/cpantesters/cron/autorun-data3.out 2>&1

# Old Website tasks

# Old Builder Entries
#1,16,36,46 0-23       * * * bash /var/www/reports/toolkit/reports-builder.sh --start >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1
#6,21,36,51 6-12,18-23 * * * bash /var/www/reports/toolkit/reports-builder-author.sh --start >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1
#3 7 * * * bash /var/www/reports/toolkit/reports-builder-author.sh --start >>/var/www/reports/toolkit/logs/builder-cron.log 2>&1

